# systemd/turtlecam.service
[Unit]
Description=TurtleCam Smart Terrarium Controller
After=network.target

[Service]
Type=simple
User=turtle
WorkingDirectory=/home/turtle/turtlecam
ExecStart=/home/turtle/turtlecam/venv/bin/python3 -m turtlecam.main
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
